# deploy/run.ini
[program:checkin]
command = %%%VENV_PATH%%%/bin/gunicorn checkin_site.wsgi
        --log-level=info
        --workers=2
        --timeout=600
        --worker-class=gevent
        --bind=0.0.0.0:65007
directory = /data/pyapp
user = root
umask = 022
priority = 999
autostart = true
startsecs = 10
startretries = 3
exitcodes = 0,2
stopsignal = TERM
stopwaitsecs = 10
stopasgroup = true
killasgroup = true
redirect_stderr = false
stdout_logfile = /data/pyapp/checkin/logs/web_stdout.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 3
stdout_capture_maxbytes = 50MB
stdout_events_enabled = false
stderr_logfile = /data/pyapp/checkin/logs/web_stderr.log
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 3
stderr_capture_maxbytes = 50MB
stderr_events_enabled = false
serverurl = AUTO